# packages/providers/test-provider/scripts/crontab
0 0 * * * cd /usr/src/app && npm run start >> /var/log/cron.log 2>&1
